# deploy/supervisor/mock_vendor.conf
[program:mock_vendor]
command=/usr/local/bin/mock-vendor-gunicorn
user=bright_tech
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/gunicorn/gunicornstdout.log
stderr_logfile=/app/logs/gunicorn/gunicornstdout.log
